# conf/meilisearch.service
[Unit]
Description=Meilisearch search engine for Bonfire
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__DATA_DIR__/meilisearch
Environment="MEILI_DB_PATH=__DATA_DIR__/meilisearch"
Environment="MEILI_HTTP_ADDR=127.0.0.1:7700"
Environment="MEILI_MASTER_KEY=__MEILI_MASTER_KEY__"
Environment="MEILI_NO_ANALYTICS=true"
Environment="MEILI_ENV=production"
ExecStart=__INSTALL_DIR__/meilisearch
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/__APP__/meilisearch.log
StandardError=inherit

# Security
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=__DATA_DIR__/meilisearch /var/log/__APP__
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
